# tools/photoshop-automation/packaging/macos/scripts/postinstall
#!/bin/zsh

set -u

LAUNCH_AGENT="/Library/LaunchAgents/com.spxad.helper.plist"
LAUNCH_AGENT_LABEL="com.spxad.helper"

/usr/sbin/chown root:wheel "$LAUNCH_AGENT"
/bin/chmod 0644 "$LAUNCH_AGENT"

CONSOLE_USER=$(/usr/bin/stat -f "%Su" /dev/console)
if [[ -z "$CONSOLE_USER" || "$CONSOLE_USER" == "root" || "$CONSOLE_USER" == "loginwindow" ]]; then
  exit 0
fi

CONSOLE_UID=$(/usr/bin/id -u "$CONSOLE_USER")

# Register the login lifecycle in the active Aqua session, then start the
# registered job through launchd so the Helper receives the normal GUI login
# environment instead of inheriting PackageKit's temporary script sandbox.
/bin/launchctl asuser "$CONSOLE_UID" /bin/launchctl bootstrap \
  "gui/$CONSOLE_UID" "$LAUNCH_AGENT" >/dev/null 2>&1 || true
/bin/launchctl asuser "$CONSOLE_UID" /bin/launchctl kickstart \
  "gui/$CONSOLE_UID/$LAUNCH_AGENT_LABEL" >/dev/null 2>&1 || true

exit 0
